# infrastructure/cluster/ansible/roles/dns-bind/files/rc.local
#!/bin/sh -e
/root/run-dns-update.sh
exit 0